# crontab
*/1 *  *   *   *     /root/cron/pinger/pinger_starter.sh
*/1 *  *   *   *     /root/cron/pinger/drower/draw_cron_v2.py > /dev/null


